# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node server.js
autorestart=true
autostart=true

[program:portfwd]
command=/usr/bin/portfwd
autorestart=true
autostart=true

